# package/debian_fh/DEBIAN/prerm
#!/bin/bash
systemctl stop corona-ad.service
systemctl disable corona-ad.service
